# etc/cron.d/botany
SHELL=/bin/sh
PATH=/usr/bin:/opt/services

0 0 * * 0	services  python /opt/services/botany/clear_weekly_users.py >/dev/null 2>&1
