# etc/cron.d/cleanup_signups
0 3 * * * root /usr/local/bin/envs.net/cleanup_signups >/dev/null 2>&1
