# etc/cron.d/cleanup_signups
0 3 * * * /usr/local/bin/envs.net/cleanup_signups
