# etc/cron.d/conntrack
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root /usr/local/bin/conntrack.sh && /usr/local/bin/envs.net/envs_conntracks.sh >/dev/null 2>&1
