# etc/cron.d/envs_sysinfo
#
# generate sysinfo.json and sysinfo.php every day
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 0 * * * root /usr/local/bin/envs/envs_sysinfo.sh >/dev/null 2>&1
