# etc/cron.d/envs_gemini
#
# generate envs gemini - index.gem (once per hour)
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 * * * * root /usr/local/bin/envs.net/envs_gemini_genpage.sh >/dev/null 2>&1
