# etc/cron.d/envs_sysinfo
#
# generate sysinfo.json and sysinfo.php once per day
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

1 0 * * * root /usr/local/bin/envs.net/envs_sysinfo.sh >/dev/null 2>&1
