# etc/cron.d/tilde_admin
#
# update tilde-admin repo (once per day)
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

0 0 * * *  root  cd /var/tilde/admins/ && git pull origin master; cp raw/banned_emails.txt /var
