# etc/cron.d/tilde_admin
#
# update tilde-admin repo (once per day)
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 0 * * *  root  /var/tilde_update_repo.sh >/dev/null 2>&1
